# jarvis-chat/scripts/jarvis-webhook-backup.service
[Unit]
Description=JARVIS Chat Backup Webhook Server
Documentation=https://github.com/MADPANDA3D/J.A.R.V.I.S-PROJECT
After=network.target docker.service jarvis-webhook.service
Requires=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT/jarvis-chat
ExecStart=/usr/bin/node /root/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT/jarvis-chat/scripts/vps-webhook-backup-server.cjs
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvis-webhook-backup

# Environment variables
Environment=NODE_ENV=production
Environment=BACKUP_WEBHOOK_PORT=9002
Environment=WEBHOOK_SECRET=your-webhook-secret-replace-this
Environment=PROJECT_ROOT=/root/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT/jarvis-chat
Environment=PRIMARY_WEBHOOK_URL=http://localhost:9000
Environment=WEBSOCKET_PORT=9004

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/root/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT/jarvis-chat

[Install]
WantedBy=multi-user.target