# jarvis-chat/scripts/jarvis-logs.service
[Unit]
Description=JARVIS Chat Log Return Service
Documentation=https://github.com/madpanda3d/jarvis-chat
After=network.target jarvis-webhook.service
Wants=jarvis-webhook.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT/jarvis-chat
ExecStart=/usr/bin/node scripts/log-return-service.js
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvis-logs

# Environment variables
Environment=NODE_ENV=production
Environment=GITHUB_TOKEN=your-github-token-here
Environment=GITHUB_REPOSITORY=madpanda3d/jarvis-chat
Environment=LOGS_DIR=/home/user/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT/logs
Environment=LOG_UPLOAD_INTERVAL=300000
Environment=MAX_LOG_SIZE=1048576
Environment=LOG_RETENTION_DAYS=7

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/user/J.A.R.V.I.S/J.A.R.V.I.S-PROJECT

[Install]
WantedBy=multi-user.target